# etc/profile.d/emacs.sh
# emacs 自启动，使用真彩色
COLORTERM=truecolor emacs --daemon &>/dev/null &
