# etc/profile.d/emacs.sh
# 在图形环境下，emacs 自启动，使用真彩色
(COLORTERM=truecolor emacsclient -nt &>/dev/null &)
